# config/systemd/user/update_bookmarks.service
[Unit]
Description=update recent bookmarks

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /home/shacon/.local/bin/update_bookmarks.py
