# local/bin/br0/br0-auto-nat.service
[Unit]
Description=br0 Wi-Fi NAT auto switch
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/home/shacon/.local/bin/br0/br0-auto-nat.sh

[Install]
WantedBy=multi-user.target
