# etc/udev/rules.d/98-ntfsfix.rules
ACTION=="add", SUBSYSTEM=="block", ENV{DEVTYPE}=="partition", RUN+="/usr/local/bin/auto-ntfsfix.sh %E{DEVNAME}"
